# docker/sphinx/crontab
* * * * * /etc/sphinxsearch/update-index.sh
* * * * * sleep 20; /etc/sphinxsearch/update-index.sh
* * * * * sleep 40; /etc/sphinxsearch/update-index.sh
